# server/dat/shadowsocks.service
[Unit]
Description=ShadowSocks - An socks5 proxy
Wants=network-online.target
After=network-online.target
Before=kcptun.service

[Service]
Type=simple
PIDFile=/run/shadowsocks.pid
ExecStart=ssserver -c /etc/ics-gfw/config.json -d start
ExecStop=ssserver -d stop
Restart=ssserver -c /etc/ics-gfw/config.json -d restart

[Install]
WantedBy=multi-user.target kcptun.service
